# queuescreen/autostart/chrome-queuescreen.desktop
#!/usr/bin/env xdg-open
[Desktop Entry]
Version=1.0
Terminal=false
Type=Application
Name=Queuescreen
Exec=sh BIN_PATH/screen.sh
Icon=chrome-kmnjmjbjmneelmflnnldpeekjdmbdhci-Default
StartupWMClass=crx_kmnjmjbjmneelmflnnldpeekjdmbdhci
